# experiments/internetscanner/templates/dns.service.j2
[Unit]
Description=sticky elephant run script
Wants=network-online.target
After=network-online.target

[Service]
Restart=always
RestartSec=1s

ExecStart=/usr/bin/docker run --name udpot -p 5053:5053/udp -p 5053:5053/tcp -d jekil/udpot

TimeoutStartSec=300
TimeoutStopSec=60
Type=forking

[Install]
WantedBy=multi-user.target
